# profile/etc/local.d/mksetup.start
#!/bin/sh

PROFILENAME=$(cat /etc/hostname)
WORKDIR=/etc/local.d

# create 'setup-$PROFILENAME'
chmod 0755 "$WORKDIR/setup"
mv -f "$WORKDIR/setup" "$WORKDIR/rootfs/usr/sbin/setup-$PROFILENAME"

# copy 'rootfs' to /
cp -rf "$WORKDIR/rootfs/." /

# copy /etc/skel content to /root
cp -r /etc/skel/. /root/

# disable 'local' service by default
rc-update del local sysinit

# this script is only used once, self-delete after usage
rm /etc/local.d/mksetup.start
